# root/etc/services.d/sonarr/run
#!/bin/sh
# Sonarr s6 service

# Required environment for .NET on FreeBSD
export HOME=/config
export XDG_CONFIG_HOME=/config/.config
export XDG_DATA_HOME=/config/.local/share
export LD_LIBRARY_PATH=/usr/local/lib
# Detect OpenSSL version (3.0=30, 3.5=35)
SSL_VER=$(openssl version | sed 's/OpenSSL \([0-9]*\)\.\([0-9]*\).*/\1\2/')
export CLR_OPENSSL_VERSION_OVERRIDE=$SSL_VER
export DOTNET_SYSTEM_NET_DISABLEIPV6=1
export System_Globalization_Invariant=true

# Signal readiness when the port is responsive
s6-ready-when fetch -qo /dev/null http://localhost:8989/ping

cd /config
exec /usr/local/bin/s6-setuidgid bsd /usr/local/share/sonarr/Sonarr -nobrowser -data=/config
